# devflow-queue.service
[Unit]
Description=DevFlow Pro Queue Worker
After=network.target

[Service]
User=roshdy
WorkingDirectory=/home/roshdy/Work/empire/dev-flow
ExecStart=/usr/bin/php artisan queue:work database --tries=1 --timeout=300
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=devflow-queue

[Install]
WantedBy=multi-user.target
